# app-sia-bot.timer
[Unit]
Description=Run app-sia-bot every hour

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
